# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx_config]
command=/opt/nginx_config
autorestart=false
startsecs=0
priority=1

[program:nginx]
command=/usr/sbin/nginx
autorestart=true
priority=2
