# supervisord.conf
[supervisord]
nodaemon=true

[program:kibana]
command=/opt/kibana-4.0.0-linux-x64/bin/kibana
autorestart=false

[program:nginx]
command=/usr/sbin/nginx
autorestart=true
